# remnant.service
[Unit]
Description=Remnant AI Agent Framework
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/remnant
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose down
Restart=always
RestartSec=5
WatchdogSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=remnant

[Install]
WantedBy=multi-user.target
